# scripts/amnezia-traffic-collector.service
[Unit]
Description=Collect Amnezia VPN telemetry
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=AMNEZIA_CONTAINER=amnezia-awg2
Environment=AMNEZIA_INTERFACE=awg0
Environment=AMNEZIA_TRAFFIC_DIR=/var/lib/amnezia-traffic
Environment=AMNEZIA_TRAFFIC_TZ=Asia/Krasnoyarsk
Environment=AMNEZIA_PING_TARGET=1.1.1.1
Environment=AMNEZIA_ACTIVE_WINDOW_SECONDS=180
Environment=AMNEZIA_PING_COUNT=3
ExecStart=/usr/local/bin/amnezia_traffic_collector.py collect
TimeoutStartSec=30
Nice=10
